# deploy/fixture-deploy.timer
[Unit]
Description=Poll GitHub for fixture.doxaed.com changes and auto-deploy

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
